# root/etc/rc.d/rc.local
#!/bin/bash

source /usr/share/fdi/commonfunc.sh
exportKCL

# setup root password
if [[ ! -z "$KCL_FDI_ROOTPW" ]]; then
  echo "root:$KCL_FDI_ROOTPW" | chpasswd >/dev/null
fi

# set encrypted root password - overrides non-encrypted password
if [[ ! -z "$KCL_FDI_ENCRYPTED_ROOTPW" ]]; then
  echo "root:$KCL_FDI_ENCRYPTED_ROOTPW" | sed -e 's/"//g' -e "s/'//g" | chpasswd -e >/dev/null
fi

# enable ssh
if [ "$KCL_FDI_SSH" == "1" ]; then
  systemctl start sshd
fi
