# root/etc/systemd/system/nm-prepare.service
[Unit]
Description=Prepares NetworkManager configuration for primary interface
Wants=basic.target
Before=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/nm-configure primary

[Install]
WantedBy=multi-user.target
